# bird-filter.timer
[Unit]
Description=update bird filter

[Timer]
OnBootSec=1min
OnUnitActiveSec=1d
Persistent=yes

[Install]
WantedBy=basic.target
